# mongod.conf
[program:mongod]
command=mongod --quiet --logpath /data/log --bind_ip_all
